# cmd/mood_trade/docker/crontab
0 * * * * root cd /app && ./mood_trade >> /var/log/cron.log 2>&1
